# packaging/debian/prerm
#!/bin/bash
# Pre-removal script

set -e

echo "🛑 Stopping OpenAuto services..."

# Stop services
if systemctl is-active --quiet openauto 2>/dev/null; then
    systemctl stop openauto
fi

if systemctl is-active --quiet openauto-btservice 2>/dev/null; then
    systemctl stop openauto-btservice
fi

# Disable services
systemctl disable openauto 2>/dev/null || true
systemctl disable openauto-btservice 2>/dev/null || true

echo "✅ OpenAuto services stopped and disabled"

exit 0
