# packaging/systemd/openauto.service
#  * Project: OpenAuto
#  * This file is part of openauto project.
#  * Copyright (C) 2025 OpenCarDev Team
#  *
#  *  openauto is free software: you can redistribute it and/or modify
#  *  it under the terms of the GNU General Public License as published by
#  *  the Free Software Foundation; either version 3 of the License, or
#  *  (at your option) any later version.
#  *
#  *  openauto is distributed in the hope that it will be useful,
#  *  but WITHOUT ANY WARRANTY; without even the implied warranty of
#  *  MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the
#  *  GNU General Public License for more details.
#  *
#  *  You should have received a copy of the GNU General Public License
#  *  along with openauto. If not, see <http://www.gnu.org/licenses/>.

[Unit]
Description=OpenAuto Android Auto Head Unit
After=network.target sound.target bluetooth.target

[Service]
Type=simple
ExecStart=/usr/bin/autoapp
Restart=on-failure
RestartSec=5s
User=root
Environment="QT_QPA_PLATFORM=eglfs"
Environment="QT_QPA_EGLFS_ALWAYS_SET_MODE=1"
Environment="DISPLAY=:0"
Environment="XDG_RUNTIME_DIR=/run/user/1000"
Environment="GST_DEBUG=3,video_sink:5"
Environment="QT_QPA_EGLFS_SWAPINTERVAL=0"
Environment="QT_QPA_EGLFS_FORCEVSYNC=0"
Environment="QT_QPA_EGLFS_INTEGRATION=eglfs_kms"
Environment="QT_QPA_EGLFS_KMS_ATOMIC=1"
Environment="QT_QPA_EGLFS_LAYER_INDEX=0"
Environment="QT_QPA_EGLFS_NO_VSYNC=1"

# Ensure audio is accessible
SupplementaryGroups=audio video plugdev

[Install]
WantedBy=multi-user.target
